# docker/supervisord.conf
[supervisord]
nodaemon=true
# Supervisord stays as root so it can open /dev/stdout for its child-process
# dispatchers; after a setuid() in this same process, Python's
# `open('/dev/stdout', 'ab')` returns EACCES because the pipe was created by
# the original root. Each [program:*] below drops to www-data via its own
# `user=`, and the children only need to *write* to their inherited stdout
# fd (already opened by root), which doesn't require the same permission.
user=root
logfile=/app/storage/supervisord.log
logfile_maxbytes=0
loglevel=info

[program:spora-web]
command=frankenphp run --config /app/frankenphp.conf --adapter caddyfile
user=www-data
autostart=true
autorestart=true
startsecs=2
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:spora-worker]
command=php /app/bin/spora worker:run --daemon
user=www-data
autostart=true
autorestart=true
startsecs=3
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
